# scripts/helios-wifi-import.service
[Unit]
Description=Helios WiFi Import — read /boot/firmware/wifi.conf into NetworkManager
After=NetworkManager.service
Requires=NetworkManager.service
Before=helios-ble.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/helios-wifi-import.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
